# ext-sources/asterisk.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='asterisk'>
  <service name='voip/asterisk' type='service' version='0'>
    <single_instance/>
    <dependency name='fs' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/local'/>
    </dependency>
    <dependency name='network' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/milestone/multi-user-server'/>
    </dependency>
    <instance name='default' enabled='true'>
      <exec_method type='method'
                   name='start'
                   exec='/usr/sbin/asterisk'
                   timeout_seconds='60'>
      </exec_method>
      <exec_method type='method'
                   name='stop'
                   exec="/usr/sbin/asterisk -rx 'core stop now'"
                   timeout_seconds='60'>
      </exec_method>
      <template>
	<common_name>
          <loctext xml:lang='C'>Asterisk : Complete IP PBX in software</loctext>      
	</common_name>
      </template>
    </instance>
  </service>
</service_bundle>
