# ext-sources/freeradius.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
        Created by Manifold
--><service_bundle type="manifest" name="freeradius">

    <service name="site/freeradius" type="service" version="1">

        <create_default_instance enabled="false"/>
        
        <single_instance/>

        <dependency name="network" grouping="require_all" restart_on="error" type="service">
            <service_fmri value="svc:/milestone/network:default"/>
        </dependency>

        <dependency name="filesystem" grouping="require_all" restart_on="error" type="service">
            <service_fmri value="svc:/system/filesystem/local"/>
        </dependency>


        
        
        
       <!-- daemon will change the UID to "freerad" or whatever is configured in /etc/raddb/radiusd.conf user=/group=  -->
        <method_context>
            <method_credential user="root" group="root"/>
        </method_context>

        <exec_method type="method" name="start" exec="/lib/svc/method/svc-freeradius start" timeout_seconds="240"/>

        <exec_method type="method" name="stop" exec=":kill" timeout_seconds="60"/>

        <property_group name="startd" type="framework">
            <propval name='duration' type='astring' value='contract'/>
            <propval name="ignore_error" type="astring" value="core,signal"/>
        </property_group>

        <property_group name="application" type="application">
            <propval name="config_file" type="astring" value="/etc/raddb/radiusd.conf"/>
        </property_group>

	<property_group name='config' type='application'>
		<!-- to enable full debugging -->
                <!-- ##TODO## improvement: more fine grained debug setting, property names might change -->
		<propval name='debug' type='boolean' value='false' />
	</property_group>
        
        
        <stability value="Evolving"/>

        <template>
            <common_name>
                <loctext xml:lang="C">
                    FreeRADIUS Server
                </loctext>
            </common_name>
        </template>

    </service>

</service_bundle>
