# ext-sources/hadoop.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='hadoop'>
  <service name='developer/distributed/hadoop' type='service' version='0'>
    <single_instance/>
    <dependency name='fs' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/local'/>
    </dependency>
    <dependency name='network' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/milestone/multi-user-server'/>
    </dependency>
    <instance name='default' enabled='false'>
      <exec_method type='method'
                   name='start'
                   exec='/usr/share/hadoop/sbin/start-all.sh'
                   timeout_seconds='300'>
	<method_context>
          <method_credential user='hadoop' group='other' />
	</method_context>
      </exec_method>
      <exec_method type='method'
                   name='stop'
                   exec='/usr/share/hadoop/sbin/stop-all.sh'
                   timeout_seconds='300'>
	<method_context>
          <method_credential user='hadoop' group='other'/>
	</method_context>
      </exec_method>
      <template>
	<common_name>
          <loctext xml:lang='C'>Hadoop - Open-source software for reliable, scalable, distributed computing.</loctext>
	</common_name>
      </template>
    </instance>
  </service>
</service_bundle>
